# ci/process_metrics-legacy.service
# process_metrics-legacy.service — event-driven process metrics collector (eBPF)
#
# Для systemd < 246 (Astra Linux CE 2.12, Debian 10, Ubuntu 18.04, RHEL 7/8).
# Использует CAP_SYS_ADMIN вместо точечных CAP_BPF/CAP_PERFMON, которые
# не поддерживаются старыми версиями systemd.
#
# Отличия от основного process_metrics.service:
#   - CAP_SYS_ADMIN вместо CAP_BPF + CAP_PERFMON (systemd < 246 не знает эти caps)
#   - MemoryLimit вместо MemoryMax (для систем с cgroup v1)
#   - Убраны ProtectClock (>= 245), ProtectKernelLogs (>= 244),
#     ProtectHostname (>= 242) — не поддерживаются в systemd 241
#   - Убран MemoryDenyWriteExecute — конфликтует с JIT-компиляцией BPF
#     на ядрах, где CAP_SYS_ADMIN нужен для BPF
#
# Установка:
#   sudo useradd -r -s /usr/sbin/nologin process_metrics
#   sudo mkdir -p /app/process_metrics
#   sudo chown nginx:nginx /app/process_metrics
#   sudo cp build/process_metrics /app/process_metrics/
#   sudo mkdir -p /app/process_metrics
#   sudo cp examples/process_metrics.conf /app/process_metrics/
#   sudo cp ci/process_metrics-legacy.service /etc/systemd/system/process_metrics.service
#   sudo cp ci/99-process-metrics.conf /etc/sysctl.d/
#   sudo sysctl --system
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now process_metrics

[Unit]
Description=Process Metrics Collector
Documentation=file:/etc/process_metrics/process_metrics.conf
After=network.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple

# --- Непривилегированный пользователь ---
User=process_metrics
Group=process_metrics

# --- Capabilities ---
AmbientCapabilities=CAP_SYS_ADMIN CAP_SYS_PTRACE CAP_KILL CAP_DAC_READ_SEARCH
CapabilityBoundingSet=CAP_SYS_ADMIN CAP_SYS_PTRACE CAP_KILL CAP_DAC_READ_SEARCH

NoNewPrivileges=yes

# --- Подготовка (от root: префикс «+») ---
ExecStartPre=+/bin/sh -c '\
  if mountpoint -q /sys/kernel/tracing 2>/dev/null; then \
    mount -o remount,mode=755 /sys/kernel/tracing 2>/dev/null || true; \
  else \
    mount -t tracefs tracefs /sys/kernel/tracing -o mode=755 2>/dev/null || true; \
  fi; \
  if ! mountpoint -q /sys/kernel/debug 2>/dev/null; then \
    mount -t debugfs debugfs /sys/kernel/debug 2>/dev/null || true; \
  fi; \
  chmod 0755 /sys/kernel/debug 2>/dev/null || true; \
  if mountpoint -q /sys/kernel/debug/tracing 2>/dev/null; then \
    mount -o remount,mode=755 /sys/kernel/debug/tracing 2>/dev/null || true; \
  else \
    mount -t tracefs tracefs /sys/kernel/debug/tracing -o mode=755 2>/dev/null || true; \
  fi'

# --- Запуск ---
ExecStart=/usr/local/bin/process_metrics -c /etc/process_metrics/process_metrics.conf

# --- Перезагрузка конфига по SIGHUP ---
ExecReload=/bin/kill -HUP $MAINPID

# --- Остановка и рестарт ---
TimeoutStopSec=5
Restart=on-failure
RestartSec=5

# --- Ограничение ресурсов ---
# MemoryLimit — cgroup v1 аналог MemoryMax.
MemoryLimit=512M
CPUQuota=50%

# --- Hardening ---
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes
ProtectKernelModules=yes
ProtectControlGroups=no
RestrictRealtime=yes
RestrictNamespaces=yes
LockPersonality=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
PrivateNetwork=no
PrivateDevices=yes

# Доступ к файловой системе (только чтение — disk I/O отсутствует)
ReadOnlyPaths=/proc /sys/fs/cgroup
ReadOnlyPaths=-/sys/kernel/btf

# BPF требует доступ к /sys/fs/bpf и debugfs/tracefs
ReadOnlyPaths=-/sys/fs/bpf
ReadOnlyPaths=-/sys/kernel/debug
ReadOnlyPaths=-/sys/kernel/debug/tracing
ReadOnlyPaths=-/sys/kernel/tracing

# Syscall фильтр
SystemCallFilter=@system-service
SystemCallFilter=bpf perf_event_open

[Install]
WantedBy=multi-user.target
