# ci/process_metrics.service
# process_metrics.service — event-driven process metrics collector (eBPF)
#
# Для systemd >= 246 (Debian 12+, Ubuntu 22.04+, RHEL 9+).
# Использует точечные capabilities CAP_BPF + CAP_PERFMON вместо CAP_SYS_ADMIN.
# Для старых систем (Astra Linux CE 2.12, Debian 10) используйте
# process_metrics-legacy.service.
#
# Запуск от непривилегированного пользователя с минимальными capabilities.
# Не требует root — только точечные caps для BPF и чтения /proc.
#
# Установка:
#   sudo useradd -r -s /usr/sbin/nologin process_metrics
#   sudo cp build/process_metrics /usr/local/bin/
#   sudo mkdir -p /etc/process_metrics
#   sudo cp examples/process_metrics.conf /etc/process_metrics/
#   sudo cp ci/process_metrics.service /etc/systemd/system/
#   sudo cp ci/99-process-metrics.conf /etc/sysctl.d/
#   sudo sysctl --system
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now process_metrics
#
# Capabilities:
#   CAP_BPF            — загрузка eBPF программ и работа с BPF-картами
#   CAP_PERFMON        — подключение к tracepoints, kprobes, raw_tracepoints
#   CAP_SYS_PTRACE     — чтение /proc/PID/* процессов других пользователей
#   CAP_KILL           — kill(pid, 0) для проверки liveness чужих процессов
#   CAP_DAC_READ_SEARCH — чтение /sys/fs/cgroup/*, /proc/*/cgroup, tracefs
#
# Kernel requirements (см. ci/99-process-metrics.conf):
#   kernel.perf_event_paranoid <= 2  (Debian default=3 блокирует CAP_PERFMON)
#   kernel.unprivileged_bpf_disabled <= 1  (Debian default=2 блокирует CAP_BPF)
#   Kernel lockdown != confidentiality  (блокирует kprobes и BPF чтение памяти)
#
# Ресурсы (по результатам стресс-тестирования, ~400 PIDs, все коллекторы вкл,
# snapshot_interval=3s, file_tracking с широким include, BPF_FILE_PATH_MAX=1024):
#
#   CPU:  ~1.3% system-wide (0.4% userspace + 0.9% BPF ядро).
#       При стрессе (12 I/O workers, 86K file events/20s): ~1.3% total.
#       Деградация (ringbuf drops) начинается при CPUQuota < 15%.
#
#   RAM:  Userspace RSS ≈ 375 МБ (HTTP buffer 256 МБ + хеш-таблицы + stacks).
#       BPF карты ядра (НЕ в cgroup): proc_map 280 МБ, fd_map 71 МБ,
#       ring buffers 390 МБ (proc 128 + file 128 + net 128 + остальные 6),
#       sock_map/tracked/tid 26 МБ = ~770 МБ.
#       MemoryMax=1792M — достаточно с запасом для userspace RSS.
#
#   Disk: 0 (все данные в RAM, disk I/O отсутствует).

[Unit]
Description=Process Metrics Collector (eBPF)
Documentation=file:/etc/process_metrics/process_metrics.conf
After=network.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple

# --- Непривилегированный пользователь ---
User=process_metrics
Group=process_metrics

# --- Capabilities (точечные, без CAP_SYS_ADMIN) ---
AmbientCapabilities=CAP_BPF CAP_PERFMON CAP_SYS_PTRACE CAP_KILL CAP_DAC_READ_SEARCH
CapabilityBoundingSet=CAP_BPF CAP_PERFMON CAP_SYS_PTRACE CAP_KILL CAP_DAC_READ_SEARCH

NoNewPrivileges=yes

# --- Подготовка (от root: префикс «+») ---
# tracefs может быть не смонтирован (Astra Linux, hardened-ядра).
# libbpf читает из него perf event ID для tracepoints.
# Пробуем оба пути: новый /sys/kernel/tracing и legacy /sys/kernel/debug/tracing.
ExecStartPre=+/bin/sh -c '\
  mountpoint -q /sys/kernel/tracing 2>/dev/null \
  || mount -t tracefs tracefs /sys/kernel/tracing 2>/dev/null \
  || true; \
  mountpoint -q /sys/kernel/debug 2>/dev/null \
  || mount -t debugfs debugfs /sys/kernel/debug 2>/dev/null \
  || true; \
  mountpoint -q /sys/kernel/debug/tracing 2>/dev/null \
  || mount -t tracefs tracefs /sys/kernel/debug/tracing 2>/dev/null \
  || true; \
  chmod 0755 /sys/kernel/tracing 2>/dev/null || true; \
  chmod 0755 /sys/kernel/debug 2>/dev/null || true; \
  chmod 0755 /sys/kernel/debug/tracing 2>/dev/null || true'

# --- Запуск ---
ExecStart=/usr/local/bin/process_metrics -c /etc/process_metrics/process_metrics.conf

# --- Перезагрузка конфига по SIGHUP ---
ExecReload=/bin/kill -HUP $MAINPID

# --- Остановка и рестарт ---
TimeoutStopSec=30
Restart=on-failure
RestartSec=5

# --- Ограничение ресурсов ---
MemoryMax=1792M
CPUQuota=25%

# --- Hardening ---
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes
ProtectKernelModules=yes
ProtectControlGroups=no
ProtectKernelLogs=yes
ProtectClock=yes
ProtectHostname=yes
RestrictRealtime=yes
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
PrivateNetwork=no
PrivateDevices=yes

# Доступ к файловой системе (только чтение — disk I/O отсутствует)
ReadOnlyPaths=/proc /sys/fs/cgroup /sys/kernel/btf

# BPF требует доступ к /sys/fs/bpf и debugfs/tracefs
ReadOnlyPaths=/sys/fs/bpf
ReadOnlyPaths=-/sys/kernel/debug
ReadOnlyPaths=-/sys/kernel/debug/tracing
ReadOnlyPaths=-/sys/kernel/tracing

# Syscall фильтр: разрешаем bpf() и perf_event_open()
SystemCallFilter=@system-service
SystemCallFilter=bpf perf_event_open

[Install]
WantedBy=multi-user.target
