# ci/process_metrics-legacy.service
# process_metrics-legacy.service — event-driven process metrics collector (eBPF)
#
# Для systemd < 246 (Astra Linux CE 2.12, Debian 10, Ubuntu 18.04, RHEL 7/8).
# Использует CAP_SYS_ADMIN вместо точечных CAP_BPF/CAP_PERFMON, которые
# не поддерживаются старыми версиями systemd.
#
# Отличия от основного process_metrics.service:
#   - CAP_SYS_ADMIN вместо CAP_BPF + CAP_PERFMON (systemd < 246 не знает эти caps)
#   - MemoryLimit вместо MemoryMax (для систем с cgroup v1)
#   - Убраны ProtectClock (>= 245), ProtectKernelLogs (>= 244),
#     ProtectHostname (>= 242) — не поддерживаются в systemd 241
#   - Убран MemoryDenyWriteExecute — конфликтует с JIT-компиляцией BPF
#     на ядрах, где CAP_SYS_ADMIN нужен для BPF (ядро выполняет JIT в контексте
#     процесса, W^X-политика блокирует это)
#
# Установка:
#   sudo useradd -r -s /usr/sbin/nologin process_metrics
#   sudo mkdir -p /var/lib/process_metrics
#   sudo chown process_metrics:process_metrics /var/lib/process_metrics
#   sudo cp build/process_metrics /usr/local/bin/
#   sudo mkdir -p /etc/process_metrics
#   sudo cp examples/process_metrics.conf /etc/process_metrics/
#   sudo cp ci/process_metrics-legacy.service /etc/systemd/system/process_metrics.service
#   sudo cp ci/99-process-metrics.conf /etc/sysctl.d/
#   sudo sysctl --system
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now process_metrics
#
# Capabilities:
#   CAP_SYS_ADMIN      — заменяет CAP_BPF + CAP_PERFMON (для ядер/systemd без их поддержки)
#   CAP_SYS_PTRACE     — чтение /proc/PID/* процессов других пользователей
#   CAP_KILL           — kill(pid, 0) для проверки liveness чужих процессов
#   CAP_DAC_READ_SEARCH — чтение /sys/fs/cgroup/*, /proc/*/cgroup, tracefs

[Unit]
Description=Process Metrics Collector (eBPF)
Documentation=file:/etc/process_metrics/process_metrics.conf
After=network.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple

# --- Непривилегированный пользователь ---
User=process_metrics
Group=process_metrics

# --- Capabilities ---
# CAP_SYS_ADMIN — необходим для bpf() и perf_event_open() на ядрах, где
# CAP_BPF/CAP_PERFMON либо не поддерживаются ядром (< 5.8), либо не распознаются
# данной версией systemd (< 246).
AmbientCapabilities=CAP_SYS_ADMIN CAP_SYS_PTRACE CAP_KILL CAP_DAC_READ_SEARCH
CapabilityBoundingSet=CAP_SYS_ADMIN CAP_SYS_PTRACE CAP_KILL CAP_DAC_READ_SEARCH

NoNewPrivileges=yes

# --- Подготовка (от root: префикс «+») ---
# tracefs может быть не смонтирован (Astra Linux, hardened-ядра).
# libbpf читает из него perf event ID для tracepoints.
# Пробуем оба пути: новый /sys/kernel/tracing и legacy /sys/kernel/debug/tracing.
# На Astra Linux CE 2.12 tracefs обычно доступен через /sys/kernel/debug/tracing.
ExecStartPre=+/bin/sh -c '\
  mountpoint -q /sys/kernel/tracing 2>/dev/null \
  || mount -t tracefs tracefs /sys/kernel/tracing 2>/dev/null \
  || true; \
  mountpoint -q /sys/kernel/debug 2>/dev/null \
  || mount -t debugfs debugfs /sys/kernel/debug 2>/dev/null \
  || true; \
  mountpoint -q /sys/kernel/debug/tracing 2>/dev/null \
  || mount -t tracefs tracefs /sys/kernel/debug/tracing 2>/dev/null \
  || true; \
  echo "=== tracefs diagnostics ==="; \
  echo "--- mounts ---"; \
  grep -E "tracefs|debugfs" /proc/mounts || echo "(none)"; \
  echo "--- /sys/kernel/tracing ---"; \
  ls -la /sys/kernel/tracing/events/sched/sched_process_exec/id 2>&1 || true; \
  echo "--- /sys/kernel/debug/tracing ---"; \
  ls -la /sys/kernel/debug/tracing/events/sched/sched_process_exec/id 2>&1 || true; \
  echo "--- permissions ---"; \
  ls -ld /sys/kernel/tracing 2>&1 || true; \
  ls -ld /sys/kernel/debug 2>&1 || true; \
  ls -ld /sys/kernel/debug/tracing 2>&1 || true; \
  echo "=== end diagnostics ==="'

# --- Запуск ---
ExecStart=/usr/local/bin/process_metrics \
    -c /etc/process_metrics/process_metrics.conf

# --- Перезагрузка конфига по SIGHUP ---
ExecReload=/bin/kill -HUP $MAINPID

# --- Рестарт при падении ---
Restart=on-failure
RestartSec=5

# --- Ограничение ресурсов ---
# BPF maps + event buffer + working set ≈ 100-120 MB with all collectors enabled.
# MemoryLimit — cgroup v1 аналог MemoryMax (для систем без cgroup v2).
# Если система на cgroup v2, используйте MemoryMax=200M вместо MemoryLimit.
MemoryLimit=200M
CPUQuota=15%

# --- Hardening ---
# Совместимо с systemd >= 231 (Debian 9+).
# ProtectClock, ProtectKernelLogs, ProtectHostname убраны — требуют systemd >= 242-245.
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes
ProtectKernelModules=yes
ProtectControlGroups=no
RestrictRealtime=yes
RestrictNamespaces=yes
LockPersonality=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
PrivateNetwork=no
PrivateDevices=yes

# Доступ к файловой системе
ReadOnlyPaths=/proc /sys/fs/cgroup
ReadOnlyPaths=-/sys/kernel/btf
ReadWritePaths=/var/lib/process_metrics /etc/process_metrics

# BPF требует доступ к /sys/fs/bpf и debugfs/tracefs
# libbpf ищет tracepoint ID в tracefs. На новых ядрах — /sys/kernel/tracing,
# на старых (Astra Linux, Debian с ядром < 4.17) — /sys/kernel/debug/tracing.
ReadOnlyPaths=-/sys/fs/bpf
ReadOnlyPaths=-/sys/kernel/debug
ReadOnlyPaths=-/sys/kernel/debug/tracing
ReadOnlyPaths=-/sys/kernel/tracing

# Syscall фильтр: разрешаем bpf() и perf_event_open()
SystemCallFilter=@system-service
SystemCallFilter=bpf perf_event_open

[Install]
WantedBy=multi-user.target
