# ci/process_metrics.service
# process_metrics.service — event-driven process metrics collector (eBPF)
#
# Для systemd >= 246 (Debian 12+, Ubuntu 22.04+, RHEL 9+).
# Использует точечные capabilities CAP_BPF + CAP_PERFMON вместо CAP_SYS_ADMIN.
# Для старых систем (Astra Linux CE 2.12, Debian 10) используйте
# process_metrics-legacy.service.
#
# Запуск от непривилегированного пользователя с минимальными capabilities.
# Не требует root — только точечные caps для BPF и чтения /proc.
#
# Установка:
#   sudo useradd -r -s /usr/sbin/nologin process_metrics
#   sudo mkdir -p /app/process_metrics
#   sudo chown process_metrics:process_metrics /app/process_metrics
#   sudo cp build/process_metrics /app/process_metrics/process_metrics
#   sudo mkdir -p /app/process_metrics
#   sudo cp examples/process_metrics.conf /app/process_metrics/
#   sudo cp ci/process_metrics.service /etc/systemd/system/
#   sudo cp ci/99-process-metrics.conf /etc/sysctl.d/
#   sudo sysctl --system
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now process_metrics
#
# Capabilities:
#   CAP_BPF        — загрузка eBPF программ и работа с BPF-картами
#   CAP_PERFMON    — подключение к tracepoints, kprobes, raw_tracepoints
#   CAP_SYS_PTRACE — чтение /proc/PID/* процессов других пользователей
#   CAP_KILL       — kill(pid, 0) для проверки liveness чужих процессов
#   CAP_DAC_READ_SEARCH — чтение /sys/fs/cgroup/*, /proc/*/cgroup, tracefs
#
# Kernel requirements (см. ci/99-process-metrics.conf):
#   kernel.perf_event_paranoid <= 2  (Debian default=3 блокирует CAP_PERFMON)
#   kernel.unprivileged_bpf_disabled <= 1  (Debian default=2 блокирует CAP_BPF)
#   Kernel lockdown != confidentiality  (блокирует kprobes и BPF чтение памяти)

[Unit]
Description=Process Metrics Collector (eBPF)
Documentation=file:/app/process_metrics/process_metrics.conf
After=network.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple

# --- Непривилегированный пользователь ---
User=process_metrics
Group=process_metrics

# --- Capabilities (точечные, без CAP_SYS_ADMIN) ---
AmbientCapabilities=CAP_BPF CAP_PERFMON CAP_SYS_PTRACE CAP_KILL CAP_DAC_READ_SEARCH
CapabilityBoundingSet=CAP_BPF CAP_PERFMON CAP_SYS_PTRACE CAP_KILL CAP_DAC_READ_SEARCH

NoNewPrivileges=yes

# --- Подготовка (от root: префикс «+») ---
# tracefs может быть не смонтирован (Astra Linux, hardened-ядра).
# libbpf читает из него perf event ID для tracepoints.
# Пробуем оба пути: новый /sys/kernel/tracing и legacy /sys/kernel/debug/tracing.
ExecStartPre=+/bin/sh -c '\
  mountpoint -q /sys/kernel/tracing 2>/dev/null \
  || mount -t tracefs tracefs /sys/kernel/tracing 2>/dev/null \
  || true; \
  mountpoint -q /sys/kernel/debug 2>/dev/null \
  || mount -t debugfs debugfs /sys/kernel/debug 2>/dev/null \
  || true; \
  mountpoint -q /sys/kernel/debug/tracing 2>/dev/null \
  || mount -t tracefs tracefs /sys/kernel/debug/tracing 2>/dev/null \
  || true; \
  chmod 0755 /sys/kernel/tracing 2>/dev/null || true; \
  chmod 0755 /sys/kernel/debug 2>/dev/null || true; \
  chmod 0755 /sys/kernel/debug/tracing 2>/dev/null || true'

# --- Запуск ---
ExecStart=/app/process_metrics/process_metrics -c /app/process_metrics/process_metrics.conf

# --- Перезагрузка конфига по SIGHUP ---
ExecReload=/bin/kill -HUP $MAINPID

# --- Рестарт при падении ---
Restart=on-failure
RestartSec=5

# --- Ограничение ресурсов ---
# BPF maps + event buffer + working set ≈ 100-120 MB with all collectors enabled.
MemoryMax=200M
CPUQuota=15%

# --- Hardening ---
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes
ProtectKernelModules=yes
ProtectControlGroups=no
ProtectKernelLogs=yes
ProtectClock=yes
ProtectHostname=yes
RestrictRealtime=yes
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
PrivateNetwork=no
PrivateDevices=yes

# Доступ к файловой системе
ReadOnlyPaths=/proc /sys/fs/cgroup /sys/kernel/btf
ReadWritePaths=/app/process_metrics

# BPF требует доступ к /sys/fs/bpf и debugfs/tracefs
# libbpf ищет tracepoint ID в tracefs. На новых ядрах — /sys/kernel/tracing,
# на старых — /sys/kernel/debug/tracing через debugfs.
ReadOnlyPaths=/sys/fs/bpf
ReadOnlyPaths=-/sys/kernel/debug
ReadOnlyPaths=-/sys/kernel/debug/tracing
ReadOnlyPaths=-/sys/kernel/tracing

# Syscall фильтр: разрешаем bpf() и perf_event_open()
SystemCallFilter=@system-service
SystemCallFilter=bpf perf_event_open

[Install]
WantedBy=multi-user.target
